# ostserver.service
[Unit]
Description=OST — Observatoire Sans Tête
After=network.target

[Service]
Environment=QT_QPA_PLATFORM=offscreen
WorkingDirectory=/home/ostserver/.ost
ExecStart=/bin/sh -c 'grep -v "^[[:space:]]*#" "$HOME/.ost/ostserver.parms" 2>/dev/null | grep -v "^[[:space:]]*$" | xargs -d "\n" /usr/bin/ostserver'
Restart=on-failure
RestartSec=5
User=ostserver

[Install]
WantedBy=multi-user.target
