# ostserver.service
[Unit]
Description=OST — Observatoire Sans Tête
After=network.target

[Service]
Environment=QT_QPA_PLATFORM=offscreen
WorkingDirectory=/home/ostserver/.ost
ExecStart=/bin/sh -c 'exec /usr/bin/ostserver $(cat "$HOME/.ost/ostserver.parms" 2>/dev/null)'
Restart=on-failure
RestartSec=5
User=ostserver

[Install]
WantedBy=multi-user.target
