# ostserver.service
[Unit]
Description=OST — Observatoire Sans Tête
After=network.target

[Service]
Environment=QT_QPA_PLATFORM=offscreen
ExecStart=/usr/bin/ostserver --webroot=/var/lib/osterix/media
Restart=on-failure
RestartSec=5
User=ostserver

[Install]
WantedBy=multi-user.target
