# cmake/cpack/postinst
#!/bin/sh
set -e

case "$1" in
    configure)
        # Use the user who invoked sudo, fallback to current user
        INSTALL_USER="${SUDO_USER:-${USER:-root}}"
        USER_HOME=$(getent passwd "${INSTALL_USER}" | cut -d: -f6)
        OST_DIR="${USER_HOME}/.ost"

        # Patch the service file with the actual user and working directory
        sed -i "s/^User=.*/User=${INSTALL_USER}/" /usr/lib/systemd/system/ostserver.service
        sed -i "s|^WorkingDirectory=.*|WorkingDirectory=${OST_DIR}|" /usr/lib/systemd/system/ostserver.service

        # Create working directory and set ownership
        mkdir -p "${OST_DIR}"
        chown "${INSTALL_USER}" "${OST_DIR}"

        # Create media directory and set ownership
        mkdir -p /var/lib/osterix/media
        chown -R "${INSTALL_USER}" /var/lib/osterix

        # Create default parms file if it does not exist yet
        PARMS_FILE="${OST_DIR}/ostserver.parms"
        if [ ! -f "${PARMS_FILE}" ]; then
            cat > "${PARMS_FILE}" << 'EOF'
# ostserver startup options — one argument per line
# This file is read at service startup and will not be overwritten on upgrade.
# Remove the leading # to activate an option.

# Web server media folder (must be writable)
--webroot=/var/lib/osterix/media

# Start embedded INDI server (Y/N)
--indiserver=Y

# Log level: 0=debug, 1=info, 2=warning, 3=error, 4=critical
--loglevel=1

# Log file path (default: ostserver_YYYY-MM-DD_hh-mm-ss.log in working directory)
#--logfile=/home/user/.ost/ostserver.log

# Security level: 0=disabled, 1=read only, 2=user grants
--grant=0

# Language (fr/en/de/it/es)
#--lng=fr

# Frontend banner text
#--banner=Observatoire Sans tête

# Load a specific configuration name (default: default)
#--configuration=default

# Modules library path (default: system path)
#--libpath=

# Database path (default: working directory)
#--dbpath=

# System resource watch interval in seconds (0 to disable)
#--systemwatchinterval=10

# Enable SSL for WebSocket server (Y/N)
#--ssl=N
#--sslcert=/etc/ost/server.crt
#--sslkey=/etc/ost/server.key

# Set admin password (applied at startup, then remove this line)
#--setadminpassword=
EOF
            chown "${INSTALL_USER}" "${PARMS_FILE}"
        fi

        # Enable and start the service (CPack builds only — debhelper uses #DEBHELPER#)
        if [ -d /run/systemd/system ]; then
            systemctl daemon-reload
            systemctl enable ostserver.service
            systemctl start ostserver.service || true
        fi
    ;;
esac

#DEBHELPER#
