# cmake/cpack/prerm
#!/bin/sh
set -e

if [ -d /run/systemd/system ]; then
    systemctl stop ostserver.service || true
    systemctl disable ostserver.service || true
    systemctl daemon-reload
fi
